# .crontab
# minute hour day_of_month month day_of_week command_to_execute

@reboot             /usr/bin/crontab > $HOME/.crontab

0   *   *   *   *   $HOME/Org/commit.sh
0   *   *   *   *   $HOME/.config/yadm/scripts/commit-and-push.sh
