# .crontab
# minute hour day_of_month month day_of_week command_to_execute

0     *     *     *     *     $HOME/.config/yadm/scripts/org-commit.sh
0     0/2   *     *     *     $HOME/.config/yadm/scripts/commit-and-push.sh
